# dist/hptdrv.conf
# hptdrv - load hpt drivers
#
# This task loads the  hpt drivers specified in the /usr/share/hptdrv directory.

description	"load hptdrvs"

start on virtual-filesystems
task
script
    drvroot=/usr/share/hptdrv
    for mod in `ls ${drvroot}`;do
      mod=`basename "${mod}"`
      if test -d "${drvroot}/${mod}"; then
        modprobe ${mod}
      fi
    done
end script
